# supervisord.conf
[supervisord]
nodaemon=true

[program:apt-cacher-ng]
command=/usr/sbin/apt-cacher-ng -c /etc/apt-cacher-ng foreground=1

[program:anacron]
command=/usr/sbin/anacron -d
